# agents/rag/com.aplan.sync-librarian-docs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Incremental sync of Datadog docs into the Librarian RAG agent.
  Runs every Sunday at 03:00 local time.

  Install:
    cp com.aplan.sync-librarian-docs.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.aplan.sync-librarian-docs.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.aplan.sync-librarian-docs.plist
    rm ~/Library/LaunchAgents/com.aplan.sync-librarian-docs.plist

  Manual trigger (without waiting for schedule):
    launchctl start com.aplan.sync-librarian-docs

  Logs:
    tail -f /Users/steve.bresnick/Projects/a_plan/agents/rag/data/librarian/sync_cron.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aplan.sync-librarian-docs</string>

    <key>WorkingDirectory</key>
    <string>/Users/steve.bresnick/Projects/a_plan/agents/rag</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/steve.bresnick/Projects/a_plan/agents/rag/.venv/bin/python</string>
        <string>sync_docs.py</string>
        <string>--agent</string>
        <string>librarian</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/steve.bresnick/Projects/a_plan/agents/rag/data/librarian/sync_cron.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/steve.bresnick/Projects/a_plan/agents/rag/data/librarian/sync_cron.log</string>
</dict>
</plist>
